# debian/xbmc-callback-daemon.upstart
## XBMC callback daemon

description "XBMC callback daemon (performs actions based on XBMC events)"
author "Peter Fern"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env DEFAULTFILE=/etc/default/xbmc-callback-daemon
env CONFIGFILE=/etc/xbmc-callback-daemon.json

script
	if [ -f "$DEFAULTFILE" ]; then
	    . "$DEFAULTFILE"
	fi

	exec xbmc-callback-daemon $CONFIGFILE
end script

post-stop exec sleep 1
